# systemd/radio-oracle.service
[Unit]
Description=Radio Oracle - Offline Voice Assistant
After=network.target ollama.service user@999.service
Wants=ollama.service
Requires=user@999.service

[Service]
Type=simple
User=oracle
Group=oracle
WorkingDirectory=/opt/radio-oracle
ExecStartPre=/bin/sleep 2
ExecStart=/opt/radio-oracle/.venv/bin/python -m oracle --mode hardware
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment
EnvironmentFile=-/opt/radio-oracle/.env
Environment=ORACLE_LOG_LEVEL=INFO
Environment=ORACLE_MODE=hardware
# Unbuffered stdio so loguru entries flow to journald in real time. When
# redirected to systemd, Python defaults to block buffering — logs sit
# in a 4 KB buffer until something forces a flush, which is invisible
# during steady-state operation.
Environment=PYTHONUNBUFFERED=1

# PulseAudio — connect to oracle's per-user pulse (AEC stack)
Environment=XDG_RUNTIME_DIR=/run/user/999
Environment=PULSE_SERVER=unix:/run/user/999/pulse/native

# Resource limits
MemoryMax=7G
CPUQuota=90%

# GPIO access
SupplementaryGroups=gpio audio

[Install]
WantedBy=multi-user.target
