# systemd/radio-oracle-diag.service
[Unit]
Description=Radio Oracle - Diagnostics Web UI
After=network.target ollama.service
Wants=ollama.service
# No Conflicts= anymore: the dashboard's read-only panels (activity feed,
# state, stats) are meant to run WHILE the radio runs. The audio/GPIO test
# cards may fail with device-busy in that case — the UI warns.

[Service]
Type=simple
User=oracle
Group=oracle
WorkingDirectory=/opt/radio-oracle
ExecStart=/opt/radio-oracle/.venv/bin/python -m oracle.diag --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment
EnvironmentFile=-/opt/radio-oracle/.env
Environment=ORACLE_LOG_LEVEL=INFO
# Same runtime dir as the main service so the state + activity files
# written by radio-oracle are found.
Environment=XDG_RUNTIME_DIR=/run/user/999

# Resource limits
MemoryMax=7G
CPUQuota=90%

# Hardware access (mic/speaker via audio group; LEDs via gpio; ADS1115 via i2c)
SupplementaryGroups=gpio audio i2c

[Install]
WantedBy=multi-user.target
